# deploy/systemd/betiq-api.service
[Unit]
Description=BetIQ REST API
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/betiq
Environment=PYTHONPATH=/opt/betiq
EnvironmentFile=-/opt/betiq/.env
ExecStart=/opt/betiq/venv/bin/python3.11 -m src.api.main
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
